# install/docker.service.d/http-proxy.conf
[Service]
Environment="HTTP_PROXY=socks5://185.239.70.124:52177" "HTTPS_PROXY=socks5://185.239.70.124:52177" "NO_PROXY=localhost,127.0.0.1,docker-registry.somecorporation.com"

